# example/lightdm
#!/usr/bin/env unitscript
### BEGIN INIT INFO
# Provides:          lightdm
# Should-Start:      console-screen kbd acpid dbus hal consolekit
# Required-Start:    $local_fs $remote_fs x11-common
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Light Display Manager
# Description:       Unitscript for the Light Display Manager
### END INIT INFO

user: root
start check: start

env scripts:
  - /etc/default/locale

program: exec lightdm
